# containers/systemd/spotweb/spotweb.container
[Unit]
Description=Spotweb container
Requires=spotweb-app.service
After=spotweb-app.service

[Container]
Image=localhost/spotweb-nginx:latest
AutoUpdate=local
LogDriver=none
UserNS=keep-id
Volume=${APP_PATH}:/app:ro,z
Network=proxy.network
Network=spotweb.network
ExposeHostPort=8080

[Service]
EnvironmentFile=%h/.config/containers/systemd/spotweb/config/app.env
TimeoutStartSec=60
Restart=on-failure

[Install]
WantedBy=multi-user.target
